# docker/build/php-fpm/supervisor/conf.d/default.conf
[program:handle-notifications]
command=/usr/bin/php /masterhome/backend/bin/console app:consumer:notifications
startsecs=10
autorestart=true
stdout_logfile=/var/www/html/logs/supervisor-%(program_name)s_%(process_num)02d.log
stderr_logfile=/var/www/html/logs/supervisor-%(program_name)s_%(process_num)02d-error.log
